# crontab
*/1 * * * * cd /home/v_shaydurov/scraper && ./run.sh >> run.log
